# linux/tunnel-reliability.service
[Unit]
Description=Jules Endpoint Tunnel Reliability Monitor
After=cloudflared.service
Wants=cloudflared.service
PartOf=cloudflared.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/bin/bash -c 'cd /opt/jules-endpoint && ./diagnostics/tunnel-reliability.sh monitor'
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment
Environment=LOG_LEVEL=INFO

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /var/log /var/run

# Resource limits
MemoryMax=128M
CPUQuota=10%

[Install]
WantedBy=multi-user.target